# files/profile.d/exports.sh
#!/bin/sh

# PATH
######
PATH="$PATH:$HOME/.vimpkg/bin"
PATH="$PATH:/usr/local/go/bin"
export PATH

# Editor
########
export EDITOR='vim'
export VISUAL='vim'

# Misc
######
export BAT_THEME='Monokai Extended'
