# files/profile.d/exports.sh
#!/bin/sh

export EDITOR='vim'
export VISUAL='vim'

export BAT_THEME='Monokai Extended'
export PATH="$PATH:$HOME/.vimpkg/bin"
