# files/profile.d/exports.sh
#!/bin/sh

export BAT_THEME='Monokai Extended'
export PATH="$PATH:$HOME/.vimpkg/bin"
